# deploy/systemd/nerdverse-ttyd.service
[Unit]
Description=Nerdverse public terminal (ttyd sandbox)
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=nerdverse-play
Group=nerdverse-play
Environment=NERDVERSE_INSTALL_ROOT=/opt/nerdverse-public
Environment=NERDVERSE_SESSION_ROOT=/var/lib/nerdverse/sessions
WorkingDirectory=/opt/nerdverse-public
ExecStart=/usr/local/bin/ttyd \
    -p 7681 \
    -i 127.0.0.1 \
    -b /play \
    -W \
    -t disableReconnect=true \
    -t enableZmodem=false \
    -t enableTrzsz=false \
    -t enableSixel=false \
    -t disableLeaveAlert=false \
    /opt/nerdverse-public/deploy/sandbox/nerdverse-cage.sh
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# sessions: per-tab active_db; telemetry: events.jsonl + stats.json aggregate
ReadWritePaths=/var/lib/nerdverse/sessions /var/lib/nerdverse/telemetry
ReadOnlyPaths=/opt/nerdverse-public

[Install]
WantedBy=multi-user.target